# private_Library/LaunchAgents/com.barsotti.notetaker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.barsotti.notetaker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/barsotti/Developer/Pers/notetaker/.venv/bin/notetaker</string>
    <string>serve</string>
    <string>--no-browser</string>
    <string>--google</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/barsotti/Developer/Pers/notetaker</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/barsotti/.notetaker/logs/serve.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/barsotti/.notetaker/logs/serve.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/barsotti/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
